# systemd/screen-lock.timer
[Timer]
OnStartupSec=30

[Install]
WantedBy=graphical.target
